# Library/LaunchAgents/junkw.xdg.environments.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>junkw.xdg.environments</string>
    <key>ProgramArguments</key>
    <array>
      <string>launchctl</string>
      <string>setenv</string>
      <string>XDG_CONFIG_HOME</string>
      <string>/Users/junkw/.config</string>
      <string>XDG_DATA_HOME</string>
      <string>/Users/junkw/.local/share</string>
      <string>XDG_CACHE_HOME</string>
      <string>/Users/junkw/.cache</string>
      <string>XDG_RUNTIME_DIR</string>
      <string>/private/tmp</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>LaunchOnlyOnce</key>
    <true/>
  </dict>
</plist>
